# deploy/supervisor/recharge-worker.conf
; ─────────────────────────────────────────────────────────────────────────────
; Supervisor configuration — Recharge Platform queue workers
;
; Installation (Ubuntu/Debian):
;   sudo apt install supervisor
;   sudo cp deploy/supervisor/recharge-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Useful commands:
;   sudo supervisorctl status
;   sudo supervisorctl restart recharge-worker:*
;   sudo supervisorctl tail -f recharge-worker:recharge-worker_00
; ─────────────────────────────────────────────────────────────────────────────

; ── High-priority recharge jobs ───────────────────────────────────────────────
; Processes the 'recharges' queue exclusively with 4 parallel workers.
; numprocs should be: (CPU cores / 2), minimum 2, maximum 8.
[program:recharge-worker]
process_name    = %(program_name)s_%(process_num)02d
command         = php /var/www/html/artisan queue:work redis \
                      --queue=recharges \
                      --sleep=1 \
                      --tries=1 \
                      --timeout=120 \
                      --max-jobs=500 \
                      --max-time=3600
autostart       = true
autorestart     = true
stopasgroup     = true
killasgroup     = true
numprocs        = 4
redirect_stderr = true
stdout_logfile  = /var/log/supervisor/recharge-worker.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 5
stopwaitsecs    = 130

; ── Retry / low-priority jobs ─────────────────────────────────────────────────
; Handles RetryRecharge and any default-queue jobs with 2 workers.
[program:recharge-retry-worker]
process_name    = %(program_name)s_%(process_num)02d
command         = php /var/www/html/artisan queue:work redis \
                      --queue=default \
                      --sleep=3 \
                      --tries=1 \
                      --timeout=120 \
                      --max-jobs=200 \
                      --max-time=3600
autostart       = true
autorestart     = true
stopasgroup     = true
killasgroup     = true
numprocs        = 2
redirect_stderr = true
stdout_logfile  = /var/log/supervisor/recharge-retry-worker.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups  = 3
stopwaitsecs    = 130

; ── Group ─────────────────────────────────────────────────────────────────────
[group:recharge-workers]
programs        = recharge-worker,recharge-retry-worker
priority        = 999
